# deploy/sabot.service
[Unit]
Description=SABot Telegram Bot
After=network.target

[Service]
Type=simple
User=sabot
Group=sabot
WorkingDirectory=/opt/sabot
EnvironmentFile=/opt/sabot/.env
Environment=PYTHONPATH=/opt/sabot
Environment=PYTHONUNBUFFERED=1

# Pre-start syntax check so failures are fast and visible in journalctl
ExecStartPre=/opt/sabot/.venv/bin/python -m py_compile /opt/sabot/main.py
ExecStart=/opt/sabot/.venv/bin/python /opt/sabot/main.py

Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3
KillMode=process
TimeoutSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sabot

# Security hardening (relaxed enough for SQLite WAL + audit.log in WorkingDirectory)
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/sabot
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
