# deploy/sabot.service
[Unit]
Description=SABot Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/SABot
ExecStart=/root/SABot/.venv/bin/python /root/SABot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
