# deploy/sabot.service
[Unit]
Description=SABot Telegram Bot
After=network.target

[Service]
Type=simple
User=sabot
Group=sabot
WorkingDirectory=/opt/sabot
EnvironmentFile=/opt/sabot/.env
Environment=PYTHONPATH=/opt/sabot
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/sabot/.venv/bin/python /opt/sabot/main.py
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3
KillMode=process
TimeoutSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sabot

# Security hardening
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sabot/data
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
